# debian/uwsgi-pythonX.Y.postinst.in
#!/bin/sh
# postinst script for @@uwsgi_pythonX.Y_binary@@
#
# see: dh_installdeb(1)

set -e

# man 3 confmodule
. /usr/share/debconf/confmodule

# First version of uwsgi source package which introduced @@uwsgi_binary@@
# binary package.
UWSGI_PYTHON_FIRST_VER="@@uwsgi_python_first_ver@@"

any_conffiles_to_migrate_are_exists() {
  [ -d /etc/uwsgi/@@uwsgi_pythonX.Y_binary@@ ] \
  && [ "x$(find /etc/uwsgi/@@uwsgi_pythonX.Y_binary@@ \
                -maxdepth 1 \
                -type f \
                -regextype posix-basic \
                -iregex '.*\.\(xml\|ini\|yaml\|yml\)$')" \
       != "x" ]
}

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

case "$1" in
    configure)
      if dpkg --compare-versions "$2" lt "$UWSGI_PYTHON_FIRST_VER" \
         && any_conffiles_to_migrate_are_exists ; then

        db_get @@uwsgi_pythonX.Y_binary@@/migrate-conffiles-to-@@uwsgi_python_binary@@
        if [ "x$RET" = "xtrue" ]; then
          find /etc/uwsgi/@@uwsgi_pythonX.Y_binary@@ \
            -maxdepth 1 \
            -regextype posix-basic \
            -iregex '.*\.\(xml\|ini\|yaml\|yml\)$' \
            -type f \
            -execdir cp '{}' /etc/@@uwsgi_python_binary@@/apps-available \; \
            -execdir sh -c '\
              CONFFILE="$(basename "{}")"
              ln -s \
                -t "/etc/@@uwsgi_python_binary@@/apps-enabled/" \
                "/etc/@@uwsgi_python_binary@@/apps-available/${CONFFILE}" \
              > /dev/null 2>&1 \
            ' \;  > /dev/null 2>&1

          invoke-rc.d @@uwsgi_python_binary@@ restart
        fi
      fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

if [ -e /usr/bin/dpkg-maintscript-helper ] \
   && dpkg-maintscript-helper supports rm_conffile 2>/dev/null; then

  for CONFFILE in \
    /etc/init.d/@@uwsgi_pythonX.Y_binary@@       \
    /etc/default/@@uwsgi_pythonX.Y_binary@@      \
    /etc/logrotate.d/@@uwsgi_pythonX.Y_binary@@  \
    /etc/uwsgi/@@uwsgi_pythonX.Y_binary@@/README \
    ; do
      dpkg-maintscript-helper rm_conffile "$CONFFILE" \
                                          "$UWSGI_PYTHON_FIRST_VER" -- "$@"
  done
fi

db_stop

exit 0
